# rootfs/etc/systemd/system/sssd.service.d/set_perms.conf
[Service]
ExecStartPre=+-/bin/chmod -f -R o= /etc/sssd
ExecStartPre=+-/bin/chmod -f -R o= /etc/sudoers.d