# rootfs/etc/systemd/system/fix_perms_nm.service
[Unit]
Description=Fixes nmconnection permssions

[Service]
Type=oneshot
ExecStart=/usr/bin/find /etc/NetworkManager/system-connections -type f -exec chmod go= {} +