# rootfs/etc/systemd/system/fix_perms_nm.service
[Unit]
Description=Fixes nmconnection permssions

[Service]
Type=oneshot
ExecStart=/bin/bash -c "chmod go= /etc/NetworkManager/system-connections/*"