# rootfs/etc/systemd/system/fix_perms_nm.path
[Unit]
Description=Watches for changes to nmconnections

[Path]
PathChanged=/etc/NetworkManager/system-connections

[Install]
WantedBy=multi-user.target
