# rootfs/etc/systemd/system/pull_images.path
[Unit]
Description=Watches for Flightctl updates

[Path]
PathChanged=/var/lib/flightctl/desired.json

[Install]
WantedBy=multi-user.target
